# src/scheduling.cron
30 5 * * * cd src && bash $HOME/src/docker_build.sh
